# asuka-mentra-g1.service
[Unit]
Description=Asuka Mentra G1 Service
After=network.target

[Service]
WorkingDirectory=%h/asuka-mentra-g1
ExecStart=/usr/bin/node dist/index.js
Environment=PORT=3000
Restart=on-failure
User=%i

[Install]
WantedBy=multi-user.target
